# root/etc/s6-overlay/s6-rc.d/svc-torrserver/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid abc /torrserver/TorrServer -d /torrserver/db -t /torrserver/db
