# root/etc/s6-overlay/s6-rc.d/init-permissions/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "running chown"
id
id abc

lsiown -R abc:abc \
    /torrserver/db

chmod -R g+w \
    /torrserver/db
